# luci/applications/luci-app-adbyby/root/etc/uci-defaults/luci-adbyby
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@adbyby[-1]
	add ucitrack adbyby
	set ucitrack.@adbyby[-1].init=adbyby
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
