# luci/applications/luci-app-nwan/root/etc/init.d/nwan
#!/bin/sh /etc/rc.common

START=95



start () {
 mkdir -p /tmp/nwan/
   sh /lib/nwan/route start 

   
}

stop () {

  echo "nwan stop"
}

